# uds-ble-id-init-service/btle-advise.timer
[Unit]
Description=Ensure BLE advise is up

[Timer]
AccuracySec=1min
OnUnitActiveSec=1min
